# rc.local
chia-blockchain/venv/bin/chia start farmer &
exit 0
