# rc.local
. ./activate
chia start farmer &
exit 0
